# systemd/.config/systemd/user/battery-status.service
[Unit]
Description=A job to check if battery is low

[Service]
Type=simple
ExecStart=%h/.local/bin/battery-status

[Install]
WantedBy=default.target
